# rPi controller/mars-joy.service
[Unit]
Description=MARS Xbox Controller Daemon
Documentation=https://github.com/your-repo/hexapod
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=starter
Group=starter
WorkingDirectory=/home/starter
ExecStart=/usr/bin/python3 -u "/home/starter/OneDrive/Projects/Robots/Hexapod_v2.0/rPi controller/joy_controller.py" --verbose
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=DISPLAY=:0
Environment=PYTHONUNBUFFERED=1

# Security hardening (optional)
NoNewPrivileges=false
PrivateTmp=false

# Give access to input devices
SupplementaryGroups=input

[Install]
WantedBy=multi-user.target
